# .init/openmedia_reduce.service
[Unit]
Description="Archive xml files"

[Service]
Type=simple
StandardOutput=journal
ExecStart="/home/some_user/app"

[Install]
